# deploy/koks.service
# systemd-Unit für den "Koks und Nutten"-Spielserver.
# Vorlage: __USER__, __DIR__, __NODE__ und __ADMIN_PASSWORD__ werden von
# deploy/install.sh ersetzt und die fertige Datei nach
# /etc/systemd/system/koks.service kopiert.

[Unit]
Description=Koks und Nutten Spielserver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__DIR__
ExecStart=__NODE__ server/dist/index.js
Restart=on-failure
RestartSec=3
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=ADMIN_PASSWORD=__ADMIN_PASSWORD__
# Nur mit SHUTDOWN_ENABLED=1 fährt der Shutdown-Button den Rechner wirklich
# herunter. In der Entwicklung (Variable nicht gesetzt) wird nur gespeichert.
Environment=SHUTDOWN_ENABLED=1

[Install]
WantedBy=multi-user.target
